# config/webcms.motleytech.supervisor.conf
[program:webcms_moteleytech]
command = /webserver/webcms/config/gunicorn_start_motleytech.sh       ; Command to start app
user = webuser                                             ; User to run as
stdout_logfile = /webserver/logs/gunicorn_motleytech_supervisor.log   ; Where to write log messages
redirect_stderr = true                                     ; Save stderr in the same log
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8            ; Set UTF-8 as default encoding

